# scripts/cpu-monitor.service
[Unit]
Description=CPU Monitor for yDaemon Service
After=network.target ydaemon.service
Requires=ydaemon.service

[Service]
Type=simple
User=root
ExecStart=/bin/bash /home/murdertxxth/git/ydaemon/scripts/cpu-monitor.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
