# scripts/cpu-monitor.service
[Unit]
Description=CPU Monitor for yDaemon Service
After=network.target ydaemon.service

[Service]
Type=simple
User=root
ExecStart=/bin/bash /opt/ydaemon/scripts/cpu-monitor.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Priority settings for system exhaustion scenarios
Nice=-10
CPUWeight=200
IOWeight=200
OOMScoreAdjust=-900

[Install]
WantedBy=multi-user.target
